# etc/systemd/system/laserflipper.service
[Unit]
Description=MKIDControl LaserFlipper Agent

[Install]
WantedBy=mkidcontrol.service

[Service]
Type=simple
ExecStart=/home/kids/anaconda3/envs/control/bin/laserflipperAgent.py
#WorkingDirectory=/home/kids/mkidcontrol/agents/
User=kids
RestartSec=100ms
Restart=always